# root/etc/profile.d/set-dconf-profile.sh
#!/bin/sh

if [ -z "$DCONF_PROFILE" ]; then
  export DCONF_PROFILE="ibus"
fi
